# deploy/launchd/com.user.ontario-land-use-chat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.ontario-land-use-chat</string>

  <key>WorkingDirectory</key>
  <string>/Users/amanv/Projects/ontario-land-use-chat</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>node_modules/next/dist/bin/next</string>
    <string>start</string>
  </array>

  <!--
    Invoke `next start` directly via node (not `bun run start`) so the
    process stays in the foreground and launchd can track its PID and
    restart it on crash. `bun run start` spawns a child next-server and
    then the bun parent exits (exit code 1), which confuses KeepAlive.
    The OPENCODE_GO_API_KEY is read from .env.local at runtime by Next.js,
    so it does not need to be set here. Set the port explicitly; the app
    must NOT collide with the WhatsApp bridge on 3000.
  -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>NODE_ENV</key>
    <string>production</string>
    <key>PORT</key>
    <string>3001</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/amanv/Projects/ontario-land-use-chat/.next/server.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/amanv/Projects/ontario-land-use-chat/.next/server.err</string>
</dict>
</plist>