# deploy/launchd/com.user.cloudflare-tunnel-chat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.cloudflare-tunnel-chat</string>

  <!--
    Runs deploy/scripts/start-tunnel.sh, which launches a Cloudflare QUICK
    tunnel to http://localhost:3001 and writes the assigned public
    trycloudflare.com URL to:

      ~/.hermes/ontario-land-use-chat-tunnel-url.txt

    Quick tunnels need no Cloudflare account but issue a RANDOM subdomain on
    every start; the URL file is rewritten each time a new URL is emitted.
    For a stable hostname, set up a named tunnel via `make tunnel-login` and
    switch this plist back to `cloudflared tunnel run ontario-land-use-chat`.
  -->
  <key>WorkingDirectory</key>
  <string>/Users/amanv/projects/ontario-land-use-chat</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/amanv/projects/ontario-land-use-chat/deploy/scripts/start-tunnel.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/amanv</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <!-- The script does its own logging to ~/.hermes; these capture any stray
       output if the script itself fails before exec'ing cloudflared. -->
  <key>StandardOutPath</key>
  <string>/tmp/cloudflared-chat.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/cloudflared-chat.err</string>
</dict>
</plist>
